# fullstack_package/frontend/security/systemd/log-watch.service
[Unit]
Description=Video Site Log Watch
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/video-site/security/monitor/log-watch.sh